# batnotify.timer
[Unit]
Description=Run battery notifier every 30 seconds 

[Timer]
OnBootSec=2min
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
